# bitcurator/env/etc/udev/rules.d/fstab.rules
# Force fstab options for specific storage devices
KERNEL=="sd*", RUN+="/usr/sbin/rbfstab"
KERNEL=="sg*", RUN+="/usr/sbin/rbfstab"
KERNEL=="fd*", ENV{ID_DRIVE_FLOPPY}="1", RUN+="/usr/sbin/rbfstab"
